# scripts/feynman-bot.service
[Unit]
Description=Feynman Lectures Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/mnt/d/Vibecoding/FeynmanLecture/feynman-bot
ExecStart=/mnt/d/Vibecoding/FeynmanLecture/feynman-bot/.venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=append:/mnt/d/Vibecoding/FeynmanLecture/feynman-bot/data/feynman-bot.log
StandardError=append:/mnt/d/Vibecoding/FeynmanLecture/feynman-bot/data/feynman-bot.log

# Secrets — fill in before enabling
Environment=TELEGRAM_BOT_TOKEN=
Environment=TELEGRAM_CHAT_ID=
Environment=ANTHROPIC_API_KEY=
Environment=DEEPSEEK_API_KEY=

[Install]
WantedBy=default.target
